# 07_Launch/src/exercise_4/config/qc_inspection.service
[Unit]
Description=Quality Control Inspection System - ROS2 Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=robot
Group=robot
WorkingDirectory=/home/robot/ros2_ws

# Environment setup
Environment="ROS_DOMAIN_ID=42"
Environment="ROBOT_ID=QC-SYSTEM-001"
Environment="PRODUCTION_LINE=MAIN-LINE-A"
Environment="QC_STATION=STATION-1"

# ROS2 launch command
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source /home/robot/ros2_ws/install/setup.bash && ros2 launch ce_robot_launch qc_inspection_launch.py'

# Restart policy for 24/7 operation
Restart=on-failure
RestartSec=10
StartLimitInterval=200
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qc-inspection

[Install]
WantedBy=multi-user.target
